# test_sd.service
[Unit]
Description=A systemd test service

[Service]
ExecStart=/usr/bin/test_sd

[Install]
WantedBy=multi-user.target


